# ansible/hab.service
[Unit]
Description=The Habitat Supervisor

[Service]
User=hab
Group=hab
ExecStart=/bin/hab sup run --ring mr

[Install]
WantedBy=default.target